# octoprint-filament-manager-mqtt/rootfs/etc/services.d/octoprint-filament-manager-mqtt/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "Starting OctoPrint Filament Manager MQTT..."

LOGLEVEL=$(bashio::config 'loglevel')
bashio::log.info "Using log level ${LOGLEVEL}"

MQTT_PORT=$(bashio::config 'mqtt.mqtt_port')
MQTT_ADDRESS=$(bashio::config 'mqtt.mqtt_address')
MQTT_TOPIC=$(bashio::config 'mqtt.mqtt_topic')
MQTT_USERNAME=$(bashio::config 'mqtt.mqtt_username')
MQTT_PASSWORD=$(bashio::config 'mqtt.mqtt_password')
bashio::log.info "Using server mqtt://${MQTT_ADDRESS}:${MQTT_PORT} with username ${MQTT_USERNAME} and topic ${MQTT_TOPIC}"

UPDATE_INTERVAL=$(bashio::config 'update_interval')

OCTOPRINT_ADDRESS=$(bashio::config 'octoprint.octoprint_address')
OCTOPRINT_API_KEY=$(bashio::config 'octoprint.octoprint_api_key')
OCTOPRINT_USE_SSL=$(bashio::config 'octoprint.octoprint_use_ssl')

export MQTT_PORT
export MQTT_ADDRESS
export MQTT_TOPIC
export MQTT_USERNAME
export MQTT_PASSWORD
export UPDATE_INTERVAL
export OCTOPRINT_ADDRESS
export OCTOPRINT_API_KEY
export OCTOPRINT_USE_SSL
export LOGLEVEL

exec python3 -u /usr/bin/octoprint-filament-manager-mqtt.py >&2
